# docs/projects/casino-manager/templates/cm-sync-hourly.crontab
0 * * * * /opt/cm/bin/cm-sync-hourly.sh >> /var/log/cm-sync-hourly.log 2>&1
